# docker/crontab
# Folderz aanbiedingen sync - 2x per dag 06:30 en 14:30
30 6 * * * cd /app && python3 sync_folderz.py >> /var/log/cron.log 2>&1
30 14 * * * cd /app && python3 sync_folderz.py >> /var/log/cron.log 2>&1

# Lege regel aan het eind vereist voor cron
